# launchd/com.user.neko8-coding-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.user.neko8-coding-agent</string>

  <!--
    tsx で a2a-server.ts を直接実行する。
    node のパスは環境で異なるため、導入前に必ず確認すること:
        which node            # 例: /opt/homebrew/bin/node
        readlink -f "$(which node)"
    nvm 等を使っている場合は解決済みの絶対パスに置き換える
    (LaunchDaemon はログインシェルを通さないため PATH も明示する)。
  -->
  <key>ProgramArguments</key>
  <array>
    <string>/opt/homebrew/bin/node</string>
    <string>/Users/bonji/workspace/localllm-agent-practice/node_modules/tsx/dist/cli.mjs</string>
    <string>src/a2a-server.ts</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/bonji/workspace/localllm-agent-practice</string>

  <key>UserName</key>
  <string>bonji</string>
  <key>GroupName</key>
  <string>staff</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>HOME</key>
    <string>/Users/bonji</string>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <!--
      接続先・公開URL は WorkingDirectory の .env からも読めるが、
      常駐の自己完結性を上げるためここでも明示しておく(任意)。
    -->
    <key>OPENAI_BASE_URL</key>
    <string>http://localhost:4000/v1</string>
    <key>OPENAI_API_KEY</key>
    <string>dummy</string>
    <key>MODEL_NAME</key>
    <string>gemma-smart</string>
    <key>A2A_PUBLIC_URL</key>
    <string>http://neko8.local:41241</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/bonji/Library/Logs/neko8-coding-agent.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/bonji/Library/Logs/neko8-coding-agent.error.log</string>
</dict>
</plist>
